# deploy/deenv.service
# systemd unit for the DeEnv kernel host. Install at /etc/systemd/system/deenv.service.
# Runs the published app as an unprivileged `deenv` user, auto-starts on boot, restarts on crash.
# See deploy/DEPLOY.md for the full runbook.

[Unit]
Description=DeEnv kernel host
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=deenv
Group=deenv

# The published binaries live here; the persistent data home is separate (DEENV_HOME), so an
# update can replace the binaries without touching kernel.json or any instance's data.
WorkingDirectory=/opt/deenv
Environment=DEENV_HOME=/var/lib/deenv

# Framework-dependent run via the installed .NET 9 runtime. If `which dotnet` differs, fix the path
# (e.g. /usr/bin/dotnet, or a self-contained build's own apphost: ExecStart=/opt/deenv/DeEnv).
ExecStart=/usr/bin/dotnet /opt/deenv/DeEnv.dll

# Always come back (crash or reboot); give the kernel time to stop its hosts cleanly on SIGTERM
# (Program.cs handles ProcessExit -> graceful shutdown).
Restart=always
RestartSec=5
TimeoutStopSec=15

# Light hardening (the app only needs to read /opt/deenv and write /var/lib/deenv).
NoNewPrivileges=true
ProtectSystem=strict
ReadWritePaths=/var/lib/deenv
ProtectHome=true

[Install]
WantedBy=multi-user.target
